# data/templates/fluentd/supervisor/supervisor.conf.erb
[program:data_enchilada_agent_<%=agent.conf_name%>]
#directory=/tmp/fluentd
command=<%=cmd%>
autostart=true
autorestart=true
#redirect_stderr=true
stderr_logfile=/var/log/data_enchilada_agent_<%=agent.conf_name%>.err.log
stdout_logfile=/var/log/data_enchilada_agent_<%=agent.conf_name%>.out.log
